# volttron/scripts/admin/volttron.service
/Users/sonakulkamnuanchai # TODO : Change Environment for Deploy Volttron as a system Services
# VOLTTRON SystemD unit
# Exec Start exectues the volttron platform script built by bootstrap
# This does not have to activate the virtualenv as the script will do that
[Unit]
Description=VOLTTRON Platform Service
After=network.target

[Service]
Type=simple
User=pea
Group=pea
WorkingDirectory=/home/pea/workspace/hive_os
#Uncomment and change this to specify a different VOLTTRON_HOME
Environment="VOLTTRON_HOME=/home/pea/.volttron"
ExecStart=/home/pea/workspace/hive_os/volttron/env/bin/volttron -vv
ExecStop=/home/pea/workspace/hive_os/volttron/env/bin/volttron-ctl shutdown --platform

[Install]
WantedBy=multi-user.target
